# scheduler/com.jobscout.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd agent — runs Job Scout on a daily schedule.

  DO NOT edit this file directly. Run this instead:
      bash scripts/setup.sh
  It fills in your actual paths and installs the agent automatically.

  After setup.sh has run:
    Verify:   launchctl list | grep jobscout
    Trigger:  launchctl start com.jobscout.daily
    Logs:     tail -f data/launchd.log
    Disable:  launchctl unload ~/Library/LaunchAgents/com.jobscout.daily.plist
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.jobscout.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>PLACEHOLDER_PYTHON</string>
    <string>run.py</string>
    <string>--headless</string>
  </array>

  <key>WorkingDirectory</key>
  <string>PLACEHOLDER_PROJECT_DIR</string>

  <!-- Schedule is filled in by scripts/setup.sh from config/search.yaml -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>   <integer>PLACEHOLDER_HOUR</integer>
    <key>Minute</key> <integer>PLACEHOLDER_MINUTE</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>PLACEHOLDER_PROJECT_DIR/data/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>PLACEHOLDER_PROJECT_DIR/data/launchd.err</string>

  <!-- Do NOT run immediately when loaded (only on schedule) -->
  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
